# systemd/user/ssh-agent.service
[Unit]
Description=SSH Agent

[Service]
ExecStart=/usr/bin/ssh-agent -D -a %t/ssh-agent.socket
Type=forking

[Install]
WantedBy=default.target

